# prod/brainzbot-plugins.service
[Unit]
Description=brainzbot plugins
After=network.target redis.service postgresql.service
Requires=redis.service postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/www/brainzbot/brainzbot-core
EnvironmentFile=-/srv/www/brainzbot/brainzbot-core/.env
ExecStart=/bin/bash -c \
  'source /srv/www/brainzbot/venv/bin/activate && \
   exec python manage.py run_plugins'
Restart=always
RestartSec=5s
KillSignal=SIGINT
TimeoutStopSec=15s
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
